# nagios.conf
description "Nagios"
start on filesystem
stop on runlevel S
respawn

# Run Nagios
exec /usr/sbin/nagios3 /etc/nagios3/nagios.cfg
